# systemd/ralphberry-dashboard.service
[Unit]
Description=Ralphberry Dashboard
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/ralphberry/apps/dashboard
Environment=NODE_ENV=production
Environment=PORT=3001
EnvironmentFile=/home/pi/ralphberry/.env
ExecStart=/usr/bin/node /home/pi/ralphberry/apps/dashboard/.next/standalone/apps/dashboard/server.js
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ralphberry-dashboard

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
