# systemd/ralphberry-bot.service
[Unit]
Description=Ralphberry Slack Bot
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/ralphberry
Environment=NODE_ENV=production
EnvironmentFile=/home/pi/ralphberry/.env
ExecStart=/usr/bin/node /home/pi/ralphberry/apps/slack-bot/dist/index.js
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ralphberry-bot

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
