# systemd/ralphberry-worker.service
[Unit]
Description=Ralphberry Worker - Job Executor
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/ralphberry
Environment=NODE_ENV=production
EnvironmentFile=/home/pi/ralphberry/.env
ExecStart=/usr/bin/node /home/pi/ralphberry/apps/worker/dist/index.js
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ralphberry-worker

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
